# contrib/systemd/ssh-agent-proxy.service
[Unit]
Description=ssh-agent-proxy: localhost SSH git-signing proxy (ssh-agent backend)
Documentation=https://github.com/pavlov-net/ssh-agent-proxy
After=default.target

[Service]
# Type=exec requires systemd >= 240; WSL2's systemd is recent enough.
Type=exec
ExecStart=%h/.local/bin/ssh-agent-proxy
EnvironmentFile=%h/.config/ssh-agent-proxy/env
Restart=on-failure
RestartSec=5s

# --- Memory hardening ---------------------------------------------------
# The process calls mlockall(MCL_CURRENT|MCL_FUTURE) at startup to keep
# transient key material out of swap. This needs a large enough locked-
# memory limit — "infinity" is fine for a single small Go binary.
LimitMEMLOCK=infinity
# Kill core dumps so key material can never land in a .core file. The
# process also sets RLIMIT_CORE=0 and PR_SET_DUMPABLE=0 internally, but
# declaring it here makes the intent explicit and covers the window
# before the in-process hardening runs.
LimitCORE=0

# --- Sandboxing ---------------------------------------------------------
# The proxy is stateless on disk: it reads one env file (via
# EnvironmentFile, which systemd loads for us *before* the service
# starts), dials the local ssh-agent Unix socket on every /sign and
# /publickey request, and serves a local HTTP listener. It never writes
# to disk, and it holds no private key material of its own — signing
# happens inside the ssh-agent process, not here. These hardening
# options reflect that.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Note: MemoryDenyWriteExecute is safe here because Go's runtime does
# not JIT and none of our dependencies need PROT_EXEC mmap at runtime.

[Install]
WantedBy=default.target
